# temp/imls-ppa/sources/imls-lockdown_1.0-1_arm/etc/systemd/system/reboot-1m.timer
[Unit]
Description=Reboot after activation

[Timer]
Unit=reboot-1m.service
OnActiveSec=60
RemainAfterElapse=false

[Install]
WantedBy=timers.target
