# temp/imls-ppa/sources/imls-lockdown_1.0-1_arm/etc/systemd/system/reboot-1m.service
[Unit]
Description=Reboot after activation

[Service]
Type=simple
ExecStart=/usr/sbin/reboot

[Install]
WantedBy=multi-user.target
