# imls-playbook/roles/unattended-upgrades/files/run-imls-client-pi-playbook-dev.timer
[Unit]
Description=Runs the IMLS playbook once a day (for dev purposes only).
After=network-online.target
Wants=network-online.target

[Timer]
Unit=run-imls-client-pi-playbook-dev.service
OnCalendar=hourly
OnBootSec=60
Persistent=true

[Install]
WantedBy=timers.target
